# files/puppetmaster.ini
[program:puppetmaster]
numprocs=3
process_name=%(process_num)02d_%(program_name)s
command=thin start -e development --log /var/log/puppet/%(process_num)02dmaster.log --socket /var/run/puppet/master.%(process_num)02d.sock --user puppet --group puppet --chdir /etc/puppet -R /etc/puppet/config.ru
startsecs=5
